# ansible/roles/dot/files/bob-update.timer
[Unit]
Description=Daily Neovim nightly update timer

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
